# Thunderbrid/thunderbird-backup.timer
[Unit]
Description=Thunderbird Backup Profile Timer

[Timer]
OnCalendar=Sun *-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target